# meta-opencentauri/recipes-apps/klipper/files/klipper-firmware-bed-init-d
#!/bin/sh
### BEGIN INIT INFO
# Provides:          klipper-firmware-bed
# Required-Start:    $local_fs $remote_fs $network
# Required-Stop:     $local_fs $remote_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start OpenCentauri firmware
### END INIT INFO

FIRMWARE="/lib/firmware/klipper-bed.bin"
FIRMWARE_VERSION="/lib/firmware/klipper-bed.bin.ver"
KATAPULT_DEPLOYER="/lib/firmware/katapult-deployer-bed.bin"
INSTALLED_VERSION="/etc/bed.ver"
SERIALPORT="/dev/ttyS4"
GPIO="201"

update(){
    flashtool -d $SERIALPORT -b 250000 -r
    flashtool -d $SERIALPORT -b 115200 -r
    if flashtool -d $SERIALPORT -f $FIRMWARE -b 115200; then
        cp $FIRMWARE_VERSION $INSTALLED_VERSION
        echo "Firmware flashed and version recorded."
    else
        echo "Firmware flash failed, version not updated."
        return 1
    fi
}

case "$1" in
    start)
        echo "Starting klipper-firmware-bed..."
        echo $GPIO > /sys/class/gpio/export
        echo out > /sys/class/gpio/gpio$GPIO/direction
        echo 1 > /sys/class/gpio/gpio$GPIO/value
        echo $GPIO > /sys/class/gpio/unexport
        sleep 1

        if [ -f "$INSTALLED_VERSION" ] && [ "$(cat $INSTALLED_VERSION)" == "$(cat $FIRMWARE_VERSION)" ]; then
            echo "Firmware is already up to date."
            exit 0
        else
            echo "Firmware version mismatch or not installed, flashing firmware..."
            if update; then
                exit 0
            else
                echo "Attempting bootloader upgrade..."
                echo $GPIO > /sys/class/gpio/export
                echo out > /sys/class/gpio/gpio$GPIO/direction
                echo 0 > /sys/class/gpio/gpio$GPIO/value
                sleep 2
                echo 1 > /sys/class/gpio/gpio$GPIO/value
                sleep 1
                echo $GPIO > /sys/class/gpio/unexport
                if mcu-flasher --firmware $KATAPULT_DEPLOYER --no-wait $SERIALPORT --firmware-version 6.4.9; then
                    echo "Bootloader deployer flashed successfully, clearing installed version."
                    rm -f $INSTALLED_VERSION
                    sleep 1
                    update
                else
                    echo "Cannot update bootloader."
                fi
            fi
        fi
        ;;
    stop)
        echo "Stopping klipper-firmware-bed..."
        echo $GPIO > /sys/class/gpio/export
        echo out > /sys/class/gpio/gpio$GPIO/direction
        echo 0 > /sys/class/gpio/gpio$GPIO/value
        echo $GPIO > /sys/class/gpio/unexport
        ;;
    restart|force-reload)
        $0 stop
        sleep 2
        $0 start
        ;;
    status)
        echo $GPIO > /sys/class/gpio/export
        gpio_value=$(cat /sys/class/gpio/gpio$GPIO/value)
        echo $GPIO > /sys/class/gpio/unexport
        if [ "$gpio_value" = "1" ]; then
            echo "klipper-firmware-bed is running"
        else
            echo "klipper-firmware-bed is not running"
        fi
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac

exit 0
